# sv/thumbor/run
#!/bin/bash

source /etc/envvars

exec 2>&1
exec thumbor -c /etc/thumbor.conf
